# recipes-nclawzero/agent-stack/files/hermes.container
# NOTE: Hermes dashboard is locked down via two layers:
#   1. Network=host so the container shares the host's network namespace —
#      the 127.0.0.1 bind below refers to host loopback, not container loopback.
#   2. The dashboard process binds to 127.0.0.1 only (not 0.0.0.0), so even
#      under host networking it's only reachable from on-box.
# Sibling containers (zeroclaw, openclaw) on the default Podman bridge cannot
# reach Hermes via container IP since Hermes is not on that bridge.
#
# Operators wanting LAN access must add an SSH tunnel or reverse proxy with
# TLS+auth from a separate process — never bind hermes to 0.0.0.0 here.
#
# Posture note: bigpi reference deployment runs hermes on the default bridge
# with PublishPort 8642:8642 (LAN-accessible without TLS). cixmini-on-Yocto
# is more conservative; LAN access requires an explicit user-mode reverse-
# proxy slice, not a quadlet default.

[Unit]
Description=hermes agent (NousResearch Hermes Agent OSS)
Documentation=https://github.com/NousResearch/hermes-agent
Wants=network-online.target nclawzero-load-agent-images.service
After=network-online.target nclawzero-load-agent-images.service

[Container]
Image=docker.io/nousresearch/hermes-agent@sha256:aa60e7483a6fad26eee233d2498d4f2b4223bf9d8990e3b07017f19b6ba7b6fe
# Image pinned to immutable digest. To roll forward, resolve a newer
# sha256 via `crane digest docker.io/nousresearch/hermes-agent:latest`
# and update here. AutoUpdate=registry intentionally OMITTED.
ContainerName=hermes
# Network=host: container shares host's network namespace. 127.0.0.1
# below is the host's loopback. PublishPort is NOT used under host
# networking (would be a no-op + warning).
Network=host
# Bind dashboard to host loopback only.
Exec=dashboard --host 127.0.0.1 --port 8642 --no-open
EnvironmentFile=/etc/nclawzero/agent-env
Volume=hermes-data:/opt/data:Z

[Service]
TimeoutStartSec=600
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target default.target
