# recipes-nclawzero/agent-stack/files/openclaw.container
[Unit]
Description=openclaw agent (upstream NemoClaw OSS)
Documentation=https://github.com/openclaw/openclaw
Wants=network-online.target nclawzero-load-agent-images.service
After=network-online.target nclawzero-load-agent-images.service

[Container]
Image=ghcr.io/openclaw/openclaw:main
# TODO(supply-chain): pin to immutable digest sha256:... once a signed
# release-tag promotion path exists. AutoUpdate=registry intentionally
# OMITTED to prevent unreviewed container code from running with API
# keys; operators trigger updates via `podman pull` plus
# `systemctl restart openclaw` on a maintenance window.
ContainerName=openclaw
EnvironmentFile=/etc/nclawzero/agent-env
Volume=openclaw-data:/openclaw-data:Z
Volume=/var/lib/nclawzero/openclaw-home:/home/node/.openclaw:Z
PublishPort=18789:18789

[Service]
TimeoutStartSec=600
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target default.target
